# packaging/systemd/magickeyboard-ui.service
[Unit]
Description=Magic Keyboard UI
After=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
ExecStart=%h/.local/bin/magickeyboard-ui
Restart=on-failure
RestartSec=2

# SteamOS desktop session is X11 here (DISPLAY=:0, no WAYLAND_DISPLAY)
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=%t
Environment=QT_QPA_PLATFORM=xcb

# Optional: less crash-noise if something goes wrong
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
