# packaging/systemd/magickeyboard-ui.service
[Unit]
Description=Magic Keyboard UI for Fcitx5
Documentation=https://github.com/HagglesworthMD/magic-keyboard
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/magickeyboard-ui
Restart=on-failure
RestartSec=2
Environment=QT_QPA_PLATFORM=wayland
Environment=QT_WAYLAND_DISABLE_WINDOWDECORATION=1

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=graphical-session.target
