# packaging/autostart/magickeyboard-ui.desktop
[Desktop Entry]
Type=Application
Name=Magic Keyboard UI
Comment=On-screen keyboard for KDE Plasma with swipe typing
GenericName=Virtual Keyboard
Exec=/usr/local/bin/magickeyboard-ui
Icon=input-keyboard
Categories=Utility;Accessibility;

# Autostart configuration
X-GNOME-Autostart-enabled=true
X-KDE-autostart-phase=2
X-KDE-autostart-after=panel

# OnlyShowIn=KDE: Intentional for SteamOS. The UI should only autostart in
# Desktop Mode (KDE Plasma), NOT in Game Mode (gamescope/Steam). This prevents
# the keyboard UI from consuming resources when the user is gaming. The fcitx5
# engine itself autostarts cross-DE (see ~/.config/autostart/fcitx5.desktop).
OnlyShowIn=KDE;

# Don't show in application launchers (it's a background service)
NoDisplay=true

# Restart on crash (KDE-specific, may not work everywhere)
X-KDE-autostart-condition=true
